# autopropose.service
[Unit]
Description=RNode Single Propose

[Service]
ExecStart=/opt/rchain-testnet-node/rundeck-scripts/propose
Type=oneshot
TimeoutStartSec=10min
