# autopropose.timer
[Unit]
Description=RNode Periodic Propose

[Timer]
OnUnitInactiveSec=1min
AccuracySec=1
RandomizedDelaySec=30
