# sync-diag.service
[Unit]
After=network.target

[Service]
ExecStart=/opt/rchain-testnet-node/sync-diag
